# puppet/scs/templates/mysqld/supervisor.conf.erb
[program:mysqld]
command = /usr/local/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe --defaults-file=/etc/mysqld/mysqld.ini
stdout_logfile = /var/log/supervisor/mysqld-stdout.log
stderr_logfile = /var/log/supervisor/mysqld-stderr.log
stopwaitsecs = 30
autorestart = true
