# scs/puppet/scs/templates/mysqld/supervisor-postload.conf.erb
[program:mysqld-postload]
command = /scs/scs/bin/run-supervisor-mysqld-postload
stdout_logfile = /scs/var/log/supervisord/%(program_name)s-stdout.log
stderr_logfile = /scs/var/log/supervisord/%(program_name)s-stderr.log
stopwaitsecs = 30
autorestart = false
